# netdiag.service
[Unit]
Description=NetDiag — Continuous network diagnostics daemon
Documentation=https://github.com/network-complex-diagnosis
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /path/to/netdiag.py --daemon --port 8080
Restart=on-failure
RestartSec=30
Environment=PYTHONUNBUFFERED=1

# Hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ProtectSystem=strict
ReadWritePaths=%h/.netdiag /tmp

[Install]
WantedBy=default.target
